# config/tiktok-api.service
[Unit]
Description=TikTok Video Intelligence API
After=network.target mongodb.service redis.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Bot-Rasperrypi
Environment="PATH=/home/pi/Bot-Rasperrypi/venv/bin"
ExecStart=/home/pi/Bot-Rasperrypi/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
